# Documentation/etc/udev.rules.d/85-jscal.rules
# Saitek PLC Cyborg Gold Joystick: jscal -p
SUBSYSTEM=="input", ATTRS{idVendor}=="06a3", ATTRS{idProduct}=="0002", ACTION=="add", RUN+="/home/bin/vayerx/jscal.sh"

